# share/pkgs/services/supervisor/centos8/supervisord.d/containerd.ini
[program:containerd]
command=/usr/share/one/supervisor/scripts/containerd.sh
directory=/tmp
priority=100
autorestart=true
startsecs=5
stopasgroup=true
killasgroup=true
redirect_stderr=true
stopwaitsecs=90
