# share/pkgs/services/supervisor/centos8/supervisord.d/crond.ini
[program:crond]
command=/usr/sbin/crond -n -s
directory=/tmp
priority=100
autorestart=true
startsecs=5
stopasgroup=true
killasgroup=true
redirect_stderr=true

