# share/pkgs/services/supervisor/centos8/supervisord.d/infinite-loop.ini
[program:infinite-loop]
command=/bin/sleep infinity
autorestart=true
stopasgroup=true
killasgroup=true
directory=/tmp
redirect_stderr=true

