# share/pkgs/services/supervisor/centos8/supervisord.d/opennebula-hem.ini
[program:opennebula-hem]
command=/usr/share/one/supervisor/scripts/opennebula-hem.sh
user=oneadmin
directory=/var/lib/one
priority=310
autorestart=true
startsecs=5
stopasgroup=true
killasgroup=true
redirect_stderr=true

