# share/pkgs/services/supervisor/centos8/supervisord.d/opennebula-httpd.ini
[program:opennebula-httpd]
command=/usr/share/one/supervisor/scripts/opennebula-httpd.sh
directory=/tmp
priority=310
autorestart=true
startsecs=5
stopasgroup=true
killasgroup=true
redirect_stderr=true

