# share/pkgs/services/supervisor/centos8/supervisord.d/opennebula-fireedge.ini
[program:opennebula-fireedge]
command=/usr/share/one/supervisor/scripts/opennebula-fireedge.sh
user=oneadmin
directory=/var/lib/one
priority=310
autorestart=true
startsecs=5
stopasgroup=true
killasgroup=true
redirect_stderr=true
stopwaitsecs=15
