# share/pkgs/services/supervisor/centos8/supervisord.d/opennebula-ssh-agent.ini
[program:opennebula-ssh-agent]
command=/usr/share/one/supervisor/scripts/opennebula-ssh-agent.sh
user=oneadmin
directory=/var/lib/one
priority=201
autostart=false
autorestart=true
startsecs=5
exitcodes=0,2
stopasgroup=true
killasgroup=true
redirect_stderr=true
stopwaitsecs=90
