# share/pkgs/services/supervisor/centos8/supervisord.d/sshd.ini
[program:sshd]
command=/usr/share/one/supervisor/scripts/sshd.sh
directory=/tmp
priority=100
autorestart=true
startsecs=5
stopasgroup=true
killasgroup=true
redirect_stderr=true

