# share/pkgs/services/supervisor/centos8/supervisord.d/stunnel.ini
[program:stunnel]
command=/usr/share/one/supervisor/scripts/stunnel.sh
directory=/tmp
priority=100
autorestart=true
startsecs=5
stopasgroup=true
killasgroup=true
redirect_stderr=true
