# share/router/vmcontext
#!/sbin/runscript

start() {
    /etc/init.d/vmcontext.rb
    return 0
}
